# journal-ai.service
[Unit]
Description=AI Trading Journal (Backend + Frontend)
Documentation=https://github.com/your-repo/ai_trading_journal
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/ai_trading_journal

# 環境變數
Environment="PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NODE_ENV=production"

# 啟動/停止命令
ExecStart=/root/ai_trading_journal/journal.sh start-prod
ExecStop=/root/ai_trading_journal/journal.sh stop

# 重啟策略
Restart=always
RestartSec=10

# 資源限制
LimitNOFILE=65536

# 系統日誌
StandardOutput=journal
StandardError=journal
SyslogIdentifier=journal-ai

[Install]
WantedBy=multi-user.target
